# roles/clamd/templates/local.clamd.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
{{ ansible_managed | comment('xml') }}
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>local.clamd</string>
        <key>UserName</key>
        <string>{{ clamd_user }}</string>
        <key>GroupName</key>
        <string>{{ clamd_group }}</string>
        <key>Program</key>
        <string>{{ clamd_path }}</string>
        <key>ProgramArguments</key>
        <array>
            <string>{{ clamd_path }}</string>
            <string>--config-file={{ clamd_config_path }}</string>
            <string>--foreground</string>
        </array>
        <key>RunAtLoad</key> <true/>
        <key>KeepAlive</key> <true/>
        <key>Nice</key> <integer>10</integer>
        <key>LowPriorityIO</key> <true/>
        <key>LowPriorityBackgroundIO</key> <true/>
    </dict>
</plist>
